# launchd/com.sudorock.tools.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.sudorock.tools</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/indy/dev/tools/start.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/indy/dev/tools</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ProcessType</key>
  <string>Background</string>

  <key>StandardOutPath</key>
  <string>/Users/indy/Library/Logs/com.sudorock.tools/stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/indy/Library/Logs/com.sudorock.tools/stderr.log</string>
</dict>
</plist>
